# statuscake/statuscake.timer
[Unit]
Description=StatusCake Timer

[Timer]
OnUnitActiveSec=250s
OnBootSec=10s

[Install]
WantedBy=timers.target
