# odroid_setup/fadecandy.service
[Unit]
Description=fadecandy

[Service]
ExecStart=/usr/local/bin/fcserver /home/odroid/2squared/config.json
Restart=always

[Install]
WantedBy=multi-user.target
